# debian/linxguard/DEBIAN/prerm
#!/bin/bash
set -e

systemctl stop linxguard || true
systemctl disable linxguard || true

exit 0
